# systemd/zapret-rpi-web-lan.service
[Unit]
Description=zapret-rpi local management UI over HTTP for Ethernet LAN
After=network-online.target zapret-rpi-web.service
Wants=network-online.target

[Service]
Type=simple
User=zapret-web
Group=zapret-web
WorkingDirectory=/usr/local/lib/zapret-rpi/web
ExecStart=/usr/local/lib/zapret-rpi/web/venv/bin/uvicorn zapret_ui.main:create_app --factory --host 0.0.0.0 --port 80 --workers 1 --no-access-log
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/zapret-rpi /etc/zapret-rpi
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
